# scripts/startup/jetson_startup.service
[Unit]
Description=RAI Jetson Orin Nano Autostart Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rai
WorkingDirectory=/home/rai/ijat2026
ExecStart=/bin/bash /home/rai/ijat2026/scripts/startup/jetson_startup.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rai-jetson-startup

[Install]
WantedBy=multi-user.target
